# systemd/test.service
[Unit]
Description=test job

[Service]
Type=simple
ExecStart=/tmp/1.sh

[Install]
WantedBy=multi-user.target
